# debian/postinst
#!/bin/bash

update-rc.d smached defaults
/etc/init.d/smached force-reload
ln -fs /usr/lib/pyshared/python2.6/smache/_native.so \
       /usr/share/pyshared/smache/_native.so
