# debian/postinst
#!/bin/bash

update-rc.d smached defaults
/etc/init.d/smached force-reload
